# alram.service
[Unit]
Description=Alram service

[Service]
Type=simple
ExecStart=/usr/bin/alram

[Install]
WantedBy=default.target