# alram.service
[Unit]
Description=Alram - Alarm for RAM

[Service]
Type=simple
ExecStart=/usr/bin/alram

[Install]
WantedBy=default.target